# packaging/systemd/user/codex-voice.service
[Unit]
Description=Codex Voice desktop dictation
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/codex-voice run
Restart=on-failure
RestartSec=2
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=graphical-session.target
