# packaging/systemd/user/codex-voice-server.service
[Unit]
Description=Codex Voice local OpenAI-compatible audio server

[Service]
Type=simple
ExecStart=/usr/local/bin/codex-voice-wrapper server --no-auth
Restart=on-failure
RestartSec=2
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
